# scripts/battery_warning.service
[Unit]
Description=Battery low warning

[Service]
Type=oneshot
ExecStart=/home/alex/Repos/arch_configs/scripts/battery_warning.sh
